# hbase-master-region/hbase-images/supervisord.conf
[supervisord]
nodaemon=true

[program:start-k8s-hbase]
command=/opt/hbase/bin/start-k8s-hbase.sh
